# PaddleOCR-VL/paddleocr/paddleocr.conf
[program:paddleocr]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/paddleocr.sh
user=root
directory=/root