# PaddleOCR-VL/paddleocr-vl/paddleocr-vl.conf
[program:paddleocr-vl]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/paddleocr-vl.sh
user=root
directory=/root
stopasgroup=true
killasgroup=true
